# deploy/systemd/ebay-upload-api.service
[Unit]
Description=eBay Upload API (TypeScript)
After=network.target

[Service]
Type=simple
User=skyler
WorkingDirectory=/home/skyler/listing-dashboard
ExecStart=/usr/bin/node packages/upload-api/dist/index.js
Restart=always
RestartSec=5
Environment=NODE_ENV=production
Environment=PORT=5051
Environment=ODOO_API_KEY=__ODOO_API_KEY__
Environment=UPLOAD_API_KEY=__UPLOAD_API_KEY__
Environment=PATH=/usr/local/bin:/usr/bin:/bin

[Install]
WantedBy=multi-user.target
